# kprop.service
[Unit]
Description=Kerberos 5 Propagation
After=syslog.target network.target
AssertPathExists=/var/kerberos/krb5kdc/kpropd.acl

[Service]
Type=forking
EnvironmentFile=-/etc/sysconfig/kprop
ExecStart=/usr/sbin/kpropd $KPROPD_ARGS

[Install]
WantedBy=multi-user.target
